# systemd/kimodo-text-encoder.service
[Unit]
Description=Kimodo Text Encoder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/CODE/blenda/mocap_research
ExecStart=%h/CODE/blenda/mocap_research/scripts/run_text_encoder_service.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
